# services/generator_control/run
#!/bin/sh
echo "*** starting GeneratorRamp ***"
exec 2>&1
exec /data/EVPGeneratorRamp/generator_ramp.py