# services/GeneratorControl/run
#!/bin/sh
echo "*** starting GeneratorControl ***"
exec 2>&1
exec /data/EVP15KVAGeneratorControl/generator_control.py
