# services/GeneratorControl/run
#!/bin/sh
echo "*** starting GeneratorControl ***"
exec 2>&1
exec /data/EVPGeneratorControl/generator_control.py
